# dist/adsyncd-0.1.2/lib/systemd/system/adsyncd.service
[Unit]
Description=Azure AD Synchronization Daemon

[Service]
Type=forking
ExecStart=/usr/bin/adsync start
ExecStop=/usr/bin/adsync stop

[Install]
WantedBy=multi-user.target
